# supervisord/openresty.conf
[program:openresty]

command=/usr/local/openresty/nginx/sbin/nginx -c /usr/local/openresty/nginx/conf/nginx.conf
autostart=true
user=www
nodaemon=false